# infra/scripts/nautilus-enclave.service
[Unit]
Description=Nautilus Nitro Enclave Service
After=network.target nitro-enclaves-allocator.service
Requires=nitro-enclaves-allocator.service

[Service]
Type=forking
User=root
WorkingDirectory=/opt/nautilus
ExecStart=/opt/nautilus/start-enclave.sh
ExecStop=/usr/bin/nitro-cli terminate-enclave --all
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
